# docker/systemd/zms.target
[Unit]
Description=ZMS/Zope and all of its dependencies
After=network-online.target

[Install]
WantedBy=multi-user.target
